# baby-feeding-bot.service
[Unit]
Description=Baby Feeding Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/baby-feeding-bot
Environment="BABY_FEEDING_BOT_TOKEN=YOUR_BOT_TOKEN"
Environment="GEMINI_API_KEY=YOUR_GEMINI_API_KEY"
ExecStart=/home/ubuntu/baby-feeding-bot/venv/bin/python /home/ubuntu/baby-feeding-bot/baby_feeding_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target